# software/ros2_ws/src/rescue_bringup/systemd/c920-stream.service
[Unit]
Description=RoboCorea C920 SRT video streamer (Jetson -> workstation GUI)
# Streams the Jetson's C920 USB camera(s) to the operator GUI over SRT (H.264,
# +Opus mic on the A/V camera). NOT a ROS node and NOT on DDS/Zenoh — a plain
# GStreamer pipeline. Edit cameras/ports/bitrate at the top of the script.
# Independent of ROS, so it does not need the Zenoh router.
After=network-online.target
Wants=network-online.target

[Service]
# The streamer is a standalone bash script (no build step). Adjust the path if
# your checkout differs (%h is $HOME). It runs the pipelines in the foreground
# and cleans them up on SIGINT/SIGTERM, so systemd tears the cgroup down cleanly.
ExecStart=/bin/bash %h/RoboCorea/software/ros2_ws/src/gui/scripts/c920_srt_stream.sh
KillSignal=SIGINT
KillMode=control-group
TimeoutStopSec=10
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
