# software/ros2_ws/src/rescue_bringup/systemd/jetson-sensors.service
[Unit]
Description=RoboCorea Jetson I2C sensors (MLX90640 thermal + LIS3MDL magnetometer)
# On-demand, like the ZED/lidar units: no [Install]/WantedBy, so it does not
# auto-start at boot. robot_manager starts it via `systemctl --user start` on a
# /robot/i2c/start request. One process owns the shared I2C bus (both nodes),
# which avoids two processes contending on the bus; per-sensor activate/deactivate
# is done at runtime via /sensors/enable_mask (bit0 mag, bit1 thermal), which
# stops a disabled sensor from touching the bus at all.

[Service]
# Must match the workstation/GUI domain (and the other robot nodes).
Environment=ROS_DOMAIN_ID=20
# %h = the user's home. Adjust the workspace path if yours differs.
ExecStart=/bin/bash -lc 'source /opt/ros/humble/setup.bash && source %h/RoboCorea/software/ros2_ws/install/setup.bash && exec ros2 launch jetson_sensors jetson_sensors.launch.py'
# Graceful SIGINT first (clean node teardown), then the cgroup is SIGKILLed.
KillSignal=SIGINT
KillMode=control-group
TimeoutStopSec=10
Restart=no
