# software/ros2_ws/src/rescue_bringup/systemd/rescue-sensors.target
[Unit]
Description=RoboCorea sensor frontend (ZED + RPLidar)
# zed.service and lidar.service declare `PartOf=` + `WantedBy=` this target, so
# starting it starts both and stopping/restarting it propagates to both.
