# software/ros2_ws/src/rescue_bringup/systemd/zed.service
[Unit]
Description=RoboCorea ZED camera driver
PartOf=rescue-sensors.target
# The ZED SDK occasionally fails to open the camera on a cold start (USB
# enumeration/bandwidth timing) — historically needing a manual retry 2-3×.
# Let systemd auto-retry: up to 6 starts per 120 s, then give up (camera absent).
StartLimitIntervalSec=120
StartLimitBurst=6

[Service]
# Must match the workstation/GUI domain (and the other robot nodes).
Environment=ROS_DOMAIN_ID=20
# publish_tf:=false so the ZED never fights the robot_localization EKF for
# odom->base_footprint. camera_name:=zed pins the topic namespace to /zed/* so it
# stays the same whatever the camera_model is (ZED1 vs ZED2 only changes the
# model/calibration, not the topics). Set camera_model to match the real camera.
ExecStart=/bin/bash -lc 'source /opt/ros/humble/setup.bash && source %h/RoboCorea/software/ros2_ws/install/setup.bash && exec ros2 launch zed_wrapper zed_camera.launch.py camera_model:=zed2 camera_name:=zed publish_tf:=false publish_map_tf:=false'
# Graceful SIGINT first; if the unit ever overruns the timeout the whole cgroup
# is SIGKILLed. The ZED shuts down cleanly on SIGINT within this window, so no
# manual ExecStopPost sweep is needed (a `pkill -f` sweep self-matches its own
# command line and marked the unit 'failed' on every stop).
KillSignal=SIGINT
KillMode=control-group
TimeoutStopSec=20
# Auto-retry a failed cold start (bounded by StartLimit* above). A clean SIGINT
# stop exits 0, so stopping from the GUI does NOT trigger a restart.
Restart=on-failure
RestartSec=5

[Install]
WantedBy=rescue-sensors.target
