# software/ros2_ws/src/rescue_bringup/systemd/zed.service
[Unit]
Description=RoboCorea ZED camera driver
PartOf=rescue-sensors.target

[Service]
# Must match the workstation/GUI domain (and the other robot nodes).
Environment=ROS_DOMAIN_ID=20
# publish_tf:=false so the ZED never fights the robot_localization EKF for
# odom->base_footprint. Adjust the workspace path / camera_model if yours differ.
ExecStart=/bin/bash -lc 'source /opt/ros/humble/setup.bash && source %h/RoboCorea/software/ros2_ws/install/setup.bash && exec ros2 launch zed_wrapper zed_camera.launch.py camera_model:=zed2 publish_tf:=false publish_map_tf:=false'
# Graceful SIGINT first, then the cgroup is SIGKILLed after the timeout. The ZED
# SDK can be slow to release the camera, so give it room; the ExecStopPost sweep
# is a belt-and-suspenders kill of any component container that outlives the unit.
KillSignal=SIGINT
KillMode=control-group
TimeoutStopSec=20
ExecStopPost=/bin/bash -lc 'pkill -9 -f "component_container_isolated.*zed" || true'
Restart=no

[Install]
WantedBy=rescue-sensors.target
